# php/build/supervisord.conf
[unix_http_server]
file=/run/supervisor.sock
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
[supervisord]
user=root
loglevel=info
pidfile=/var/run/supervisord.pid
logfile=/var/log/supervisor/supervisord.log
childlogdir=/var/log/supervisor

[program:php-fpm]
command=php-fpm
autostart=true
autorestart=true
startretries=999999
startsecs=3
priority=1

[program:crond]
command=crond -f
autostart=true
autorestart=true
startretries=999999
startsecs=3
priority=2

[supervisorctl]
serverurl=unix:///run/supervisor.sock 

[include]
files = /etc/supervisor.d/*.ini
